# root/etc/services.d/fake/run
#! /bin/python

I only exist to loop test our shell checking system. 

echo $TEST > /test/testfile.txt
